# supervisor/user.conf
[program:user-startup]
environment = HOME="/root",USER="root",COMPOSER_HOME="/root/.composer"
command=sh /var/www/start.sh
directory=/var/www/
autorestart=false
startsecs=10
